# find-user-dirs.timer
[Unit]
Description = find-user-dirs

[Timer]
OnStartupSec=1sec
OnUnitActiveSec=10min

[Install]
WantedBy = timers.target
